# configs/cron.hourly/btrbk
#!/bin/sh
exec /usr/local/bin/btrbk -q run
